# deploy/debot-pair-btceth-c.service
[Unit]
Description="debot-pair-btceth-c (A/B test: ez=0.2+sl=4.0+mlr=1.0)"
After=network.target

[Service]
WorkingDirectory=/opt/debot
ExecStart=/bin/bash /opt/debot/scripts/debot_lighter.sh debot-pair-btceth-c
Restart=on-failure
RestartSec=30
StartLimitBurst=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
